# dustsensorbridge/src/dist/dustsensorbridge.service
# copy this to /etc/systemd/system
# then sudo systemctl daemon-reload
# and enable with sudo systemctl enable --now dustsensorbridge

[Unit]
Description=dustsensorbridge
After=network.target

[Service]
WorkingDirectory=/home/bertrik/dustsensorbridge
ExecStart=/home/bertrik/dustsensorbridge/bin/dustsensorbridge
Type=simple
KillMode=process
Restart=on-failure
RestartSec=10s
TimeoutStopSec=15s
User=bertrik

[Install]
WantedBy=multi-user.target
